# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

; 1. PostgreSQL — starts first
[program:postgres]
command=su-exec postgres postgres -D /var/lib/postgresql/data
autostart=true
autorestart=true
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 2. Go backend — waits for PG via retry
[program:backend]
command=/app/server
autostart=true
autorestart=true
startretries=20
startsecs=2
priority=200
; Fixed vars; secrets (GEMINI_API_KEY, JWT_SECRET, TELEGRAM_BOT_TOKEN) are
; inherited from the container environment (docker run -e / env_file).
environment=DATABASE_URL="postgres://tma_user:tma_password@localhost:5432/tma_boilerplate?sslmode=disable",PORT="8080",ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 3. nginx — serves frontend + proxies /api
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
